# config/crontabs
0 12 * * 1 /usr/local/bin/letsencrypt-renew
